# bazarr/rootfs/etc/s6-overlay/s6-rc.d/bazarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Bazarr
# Runs Bazarr
# ==============================================================================
bashio::log.info 'Starting Bazarr...'

python3 /opt/bin/bazarr.py --no-update --config /config
